# etc/crontabs/root
# Custom config

# Update opkg database every hour
0 * * * * opkg update
# Auto commit and push changes in config to git every few minutes
*/5 * * * * sh /overlay/upper/roll-changes.sh

# Added by ACME
0 0 * * * /etc/init.d/acme start
